# system_snapshot/20250930_144001/systemd-overrides/snap-snap\x2dstore-1216.mount
[Unit]
Description=Mount unit for snap-store, revision 1216
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snap-store_1216.snap
Where=/snap/snap-store/1216
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
